# config/includes.chroot/etc/skel/.config/autostart/flameshot.desktop
[Desktop Entry]
Name=Flameshot
Comment=Screenshot tool
Exec=/usr/bin/flameshot
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
